# scripting/etc/systemd/system/bydefault.service
[Unit]
Description=Run script by default
After=network.target

[Service]
Environment=DISPLAY=:0
ExecStart=/opt/bydefault/bydefault
WorkingDirectory=/opt/bydefault/
Restart=always
RestartSec=2
User=pi

[Install]
WantedBy=graphical.target

