# scripting/etc/systemd/system/create-data.service
[Unit]
Description=create data partition

[Service]
Type=oneshot
ExecStart=/boot/bydefault/create_data

[Install]
WantedBy=multi-user.target

