# bootstrap_root/etc/systemd/system/create-data.service
[Unit]
Description=create data partition

[Service]
Type=oneshot
ExecStart=/usr/bin/create_data

[Install]
WantedBy=multi-user.target

